# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:webserver]
command=node /app/dist/index.cjs
directory=/app
environment=NODE_ENV="production",DATABASE_URL="/app/data/data.db"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:engine]
command=/app/engine/.venv/bin/python engine.py
directory=/app/engine
environment=API_BASE_URL="http://localhost:3000"
autostart=true
autorestart=true
startretries=10
startsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2
